# instancesetup/crontab
*/5 * * * * /home/ubuntu/autoshutdown.sh > /tmp/cronjob.log 2>&1
